# roles/supervisor/templates/program.conf
{% set program_names = [] %}
{% for process in worker.processes %}
{# expression-statement is not supported #}
[program:{{ process.workerName }}]
process_name=%(program_name)s_%(process_num)02d
directory=/srv/{{ repo.name }}
command={{ process.command }}
{% if process.number is not number %}
{% set hv = hostvars[inventory_hostname] %}
{% set var_name = process.number %}
{% set process_number = hv[var_name] | default('0') %}
numprocs={{ process_number }}
{% else %}
numprocs={{ process.number }}
{% endif %}
{% if process.redirect_stderr is defined %}
redirect_stderr={{ process.redirect_stderr }}
{% endif %}
{% if process.autostart is defined %}
autostart={{ process.autostart }}
{% else %}
autostart=true
{% endif %}
{% if process.autorestart is defined %}
autorestart= {{ process.autorestart }}
{% else %}
autorestart=true
{% endif %}
{% if process.environment is defined %}
environment={{ process.environment }}
{% endif %}
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-stderr.log
{% endfor %}

[group:{{ repo.name }}]
programs={{ ','.join(program_names) }}
